# bookma.service
[Unit]
Description=BookMa - Beautiful Bookmark Manager
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/bookma
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000

# Logging
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=bookma

[Install]
WantedBy=multi-user.target
